# actividad4/shellscript.service
[Unit]
Description=Actividad 4 SO1

[Service]
ExecStart=/usr/bin/script.sh

[Install]
WantedBy=multi-user.target